# hermes/rootfs/etc/services.d/haproxy/run
#!/usr/bin/with-contenv sh
set -euo pipefail
exec /usr/sbin/haproxy -W -db -f /etc/haproxy/haproxy.cfg -p /run/haproxy.pid
