# hermes/rootfs/etc/s6-overlay/s6-rc.d/haproxy/run
#!/command/with-contenv bash
exec /usr/sbin/haproxy -W -db -f /etc/haproxy/haproxy.cfg -p /run/haproxy.pid
